# deploy/mini-games-server.service
[Unit]
Description=Mini Games Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mini-games-server
ExecStart=/opt/mini-games-server/mini_games_server --static-files-path /opt/mini-games-server/web-client/dist
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
